# nsq/init/nsqd.conf
author "zw"
description "nsqd"
start on startup
stop on shutdown
respawn limit 10 5  

script
   cd /root/back/nsq/nsq
   exec ./testServer
end script
